# BlackBox/etc/cron.d/rewrite_maps
# Obfuscate routes
*/10 * * * * root /opt/obfuscator/update_routes.sh > /dev/null

